# agentkit/lab.agentkit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>lab.agentkit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/matus/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>agentkit.main</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/matus/homelab/agentkit</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/matus/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/matus/.local/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/matus/Library/Logs/agentkit/agentkit.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/matus/Library/Logs/agentkit/agentkit.error.log</string>
</dict>
</plist>
