# audio/com.audio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.audio</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/matus/.local/bin/uv</string>
        <string>run</string>
        <string>app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/matus/homelab/audio</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/audio.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/audio.err.log</string>
</dict>
</plist>
