# whisper/com.whisper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whisper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/whisper-server</string>
        <string>-m</string>
        <string>/Users/matus/homelab/whisper/ggml-large-v3-turbo.bin</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>9100</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/matus/homelab/whisper</string>
    <key>StandardOutPath</key>
    <string>/tmp/whisper.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/whisper.err.log</string>
</dict>
</plist>
